# vagrant/cookbooks/threatnote/files/default/threatnote.conf
[program:threatnote]
command=python server.py
directory=/threatnote
stdout_logfile=/home/vagrant/threatnote.log
stderr_logfile=/home/vagrant/threatnote-error.log